# debian/postinst
#!/bin/sh
# Setup system user
if ! getent passwd "rpanion" >/dev/null
then
    adduser --system --no-create-home --group --quiet rpanion
fi

sudo usermod -a -G plugdev rpanion
sudo usermod -a -G dialout rpanion
sudo usermod -a -G netdev rpanion
sudo usermod -a -G gpio rpanion
sudo usermod -a -G video rpanion

# Set permissions
chown -R rpanion:rpanion /usr/share/rpanion-server
chown -R rpanion:rpanion /etc/rpanion-server
chmod +x /usr/share/rpanion-server/app/server/index.js
systemctl daemon-reload
systemctl enable rpanion-server
systemctl start rpanion-server
